# backend/brain/supervisord.conf
; Combined-container supervisord config — runs nginx + Node + Python brain in one App Runner instance.
; ALL output goes to stdout/stderr (App Runner CloudWatch) so we can debug boot failures.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

; Forward supervisord's own logging to PID 1 stdout
[supervisorctl]

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

[program:node]
command=node /app/dist/server.js
directory=/app
environment=PORT=3002
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20

[program:brain]
command=/opt/brain/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
directory=/app/brain
environment=PYTHONPATH="/app/brain",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=30
