# backend/brain/supervisord.conf
; Combined-container supervisord config — runs nginx + Node + Python brain in one App Runner instance.
; nginx handles Host-based routing on port 3001 (external):
;   api.<heru>.com       → localhost:3002 (Node Express)
;   brain-api.<heru>.com → localhost:8080 (Python cognee + uvicorn)

[supervisord]
nodaemon=true
user=root
logfile=/tmp/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/tmp/logs/nginx.stdout.log
stderr_logfile=/tmp/logs/nginx.stderr.log
priority=10

[program:node]
command=node /app/dist/server.js
directory=/app
environment=PORT=3002
autostart=true
autorestart=true
stdout_logfile=/tmp/logs/node.stdout.log
stderr_logfile=/tmp/logs/node.stderr.log
priority=20

[program:brain]
command=/opt/brain/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
directory=/app/brain
environment=PYTHONPATH="/app/brain",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
stdout_logfile=/tmp/logs/brain.stdout.log
stderr_logfile=/tmp/logs/brain.stderr.log
priority=30
